# scripts/sensor-api.service
[Unit]
Description=Sensor Data REST API
After=network.target
Wants=network.target

[Service]
Type=simple
User=mqttapp
Group=mqttapp
WorkingDirectory=/opt/mqttapp
Environment=PYTHONPATH=/opt/mqttapp

# Путь к Python в виртуальном окружении
ExecStart=/opt/mqttapp/venv/bin/python -m src.api_server

# Перезапуск при сбоях
Restart=always
RestartSec=5

# Логирование
StandardOutput=file:/opt/mqttapp/logs/api_server.log
StandardError=file:/opt/mqttapp/logs/api_server-error.log

# Переменные окружения
Environment=PYTHONUNBUFFERED=1

# Безопасность
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
